# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Detector de personas (Python / DeepFace) ─────────────────────────────────
[program:person-detector]
; python -m uvicorn: el módulo lo resuelve EXACTAMENTE el intérprete que lo
; tiene instalado — inmune a shebangs/PATH de otro Python (causa del FATAL
; "No module named 'uvicorn'" cuando convivían dos Python 3.11 en la imagen).
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8090
directory=/app/person-detector
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Dar tiempo al modelo para cargar antes de marcar como listo
startsecs=10
priority=10

; ── Usuarios service (Java / Spring Boot) ────────────────────────────────────
; MaxRAMPercentage acota el heap a un 40% del límite de memoria del
; container, para que el JVM no compita por toda la memoria disponible con
; el sidecar de Python (TensorFlow/DeepFace) que corre en el mismo container.
[program:usuarios-service]
command=java -XX:MaxRAMPercentage=40.0 -jar /app/app.jar
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Esperar a que el detector esté listo
startsecs=5
priority=20
