# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Detector de personas (Python / DeepFace) ─────────────────────────────────
[program:person-detector]
command=uvicorn person-detector.app.main:app --host 127.0.0.1 --port 8090
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Dar tiempo al modelo para cargar antes de marcar como listo
startsecs=10
priority=10

; ── Usuarios service (Java / Spring Boot) ────────────────────────────────────
[program:usuarios-service]
command=java -jar /app/app.jar
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Esperar a que el detector esté listo
startsecs=5
priority=20
